# roles/sys-ctl-hlth-msmtp/templates/sys-ctl-hlth-msmtp.service.j2
[Unit]
Description=Check msmtp liveliness
OnFailure=sys-ctl-alm-telegram.{{ SOFTWARE_NAME }}@%n.service

[Service]
Type=oneshot
ExecStart=/bin/bash {{ health_msmtp_folder }}sys-ctl-hlth-msmtp.sh
