# roles/sys-ctl-rpr-btrfs-balancer/templates/sys-ctl-rpr-btrfs-balancer.service.j2
[Unit]
Description=auto balance btrfs
OnFailure=sys-ctl-alm-compose.{{ SOFTWARE_NAME }}@%n.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'btrfs-auto-balancer 90 10'