# roles/sys-ctl-rpr-btrfs-balancer/templates/systemctl.service.j2
[Unit]
Description=auto balance btrfs
OnFailure={{ SYS_SERVICE_ON_FAILURE_COMPOSE }}

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'btrfs-auto-balancer 90 10'