# roles/svc-opt-keyboard-color/templates/systemctl.service.j2
[Unit]
Description=Keyboard Color Service
OnFailure={{ SYS_SERVICE_ON_FAILURE_COMPOSE }}

[Service]
Type=oneshot
ExecStart={{ systemctl_script_exec }} {{ vendor_and_product_id }}