# roles/sys-ctl-mtn-cert-deploy/templates/systemctl.service.j2
[Unit]
Description=TLS cert deploy to {{ lookup('container', application_id, 'directories.instance') }}
OnFailure={{ SYS_SERVICE_ON_FAILURE_COMPOSE }}

[Service]
Type=oneshot
Environment="compose_cmd={{ compose_cmd }}"
ExecStart={{ system_service_script_exec }} \
  {{ lookup('cert', application_id, 'files.cert') | dirname }} \
  {{ lookup('container', application_id, 'directories.instance') }}
