# roles/sys-ctl-mtn-cert-deploy/templates/systemctl.service.j2
[Unit]
Description=TLS cert deploy to {{ lookup('docker', application_id, want='directories.instance') }}
OnFailure={{ SYS_SERVICE_ON_FAILURE_COMPOSE }}

[Service]
Type=oneshot
Environment="compose_cmd={{ compose_cmd }}"
ExecStart={{ system_service_script_exec }} \
  {{ lookup('cert_plan', application_id, want='files.cert') | dirname }} \
  {{ lookup('docker', application_id, want='directories.instance') }}
